# limitclock.service
[Unit]
Description=LimitClock — Claude Code rate limit reset notifier
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.limitclock
ExecStart=/usr/bin/node %h/.limitclock/index.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
