# limitclock.service
[Unit]
Description=LimitClock — Claude Code rate limit reset notifier
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/sadrikov49/Desktop/LimitClock
ExecStart=/usr/bin/node /home/sadrikov49/Desktop/LimitClock/index.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
